# pi-config/pi-display.service
[Unit]
Description=Raspberry Pi Display Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/display
ExecStart=/usr/bin/python3 /home/pi/display/startup_script.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 